# sidecar/garmin/com.apex.garmin-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Copy to ~/Library/LaunchAgents/ and fix the two absolute paths if the
     repo lives elsewhere. Load with:
     launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.apex.garmin-collector.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.apex.garmin-collector</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/abdelouahedelbied/Health app/sidecar/garmin/run-when-online.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/abdelouahedelbied/garmin-collector-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/abdelouahedelbied/garmin-collector-launchd.log</string>
</dict>
</plist>
